# templates/agent/systemd/strands-pg-deploy.path.in
[Unit]
Description=Watch @AGENT@ deploy-trigger file

[Path]
PathModified=@DIR@/.deploy-trigger
Unit=@AGENT@-deploy.service

[Install]
WantedBy=multi-user.target
